# backend/deploy/supervisor/api.conf
;# ============================================================
;# Supervisor 配置 - 宝妈英语早操后端
;# 路径: /etc/supervisor/conf.d/mom-english-backend.conf
;# ============================================================

[program:mom-english-backend]
command=/opt/mom-english-backend/venv/bin/gunicorn mom_english_backend.wsgi:application
    --bind 127.0.0.1:8000
    --workers 3
    --worker-class sync
    --timeout 60
    --max-requests 1000
    --max-requests-jitter 100
    --access-logfile /opt/mom-english-backend/logs/access.log
    --error-logfile /opt/mom-english-backend/logs/error.log
    --log-level info
directory=/opt/mom-english-backend
user=momapp
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
priority=999
environment=
    PATH="/opt/mom-english-backend/venv/bin",
    DJANGO_SETTINGS_MODULE="mom_english_backend.settings",
    PYTHONUNBUFFERED="1"

[group:mom-english]
programs=mom-english-backend
priority=999
